# supervisor.conf
[program:caronte]
command=/usr/bin/env python -u caronte.py
directory=/opt/caronte
user=root
autostart=true
autorestart=true
redirect_stderr=True
